# engine/systemd/nspawn-vault-check.service
[Unit]
Description=nspawn-vault dead-man switch: check for stale backups
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/nspawn-vault/check-stale.sh
StandardOutput=journal
StandardError=journal
